# api/debian/table-tv.service
[Unit]
Description=Table TV - pool table camera streaming and overlay
After=network.target avahi-daemon.service
Wants=avahi-daemon.service

[Service]
Type=simple
Environment=UI_DIST_PATH=/usr/share/table-tv/ui-dist
Environment=SQLITE_PATH=/var/lib/table-tv/table-tv.db
Environment=PORT=8080
ExecStart=/usr/lib/table-tv/table-tv-start.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
